# src/ingest_core/myemsl-client.cron
#!/bin/bash
yum -y install myemsl-client myemsl-common
