# src/metadata/50-fmds-versionclearer
#!/bin/bash
#
# Init file for fmds server daemons
#
# chkconfig: 2345 55 25
# description: fmds server daemon
#
# pidfile: /var/run/pacifica_metadata_fmds_versionclearer.pid

# source function library

CHANNEL=fmds
PID_FILE=/var/run/pacifica_metadata_fmds_versionclearer.pid
LOCKFILE=/var/lock/subsys/pacifica_metadata_fmds_versionclearer

. /usr/libexec/pacifica/notification/versionclearer_init
